# dotfiles/.config/autostart/skype.desktop
[Desktop Entry]
Type=Application
Name=Skype
Exec=/usr/bin/skype
